# ConfMon/arkos-confmon.service
[Unit]
Description=Update arkOS low-level configuration
 
[Service]
Type=oneshot
ExecStart=/usr/bin/confmon.py
 
[Install]
WantedBy=multi-user.target